# src/deb/control/postinst
#!/bin/bash
set -e
if ! id bayeos > /dev/null 2>&1 
then	
	useradd bayeos
fi

systemctl start postgresql.service 

chmod +x /usr/share/bayeos-server/*.sh
su - postgres -c /usr/share/bayeos-server/pg_postinst.sh

chown -R bayeos:bayeos /var/lib/bayeos-server 

chmod 700 /var/lib/bayeos-server/bayeos-server.jar
ln -sf /etc/bayeos-server/logback.xml /var/lib/bayeos-server/logback.xml

kf=/etc/bayeos-server/api_key
if [ ! -f $kf ]; then    
	tr -dc A-Za-z0-9 </dev/urandom | head -c 32 > $kf
	chown bayeos:bayeos $kf
	chmod 600 $kf
    echo "Generated new api key in $kf"
fi

systemctl enable bayeos-server.service
systemctl start bayeos-server.service  

a2enmod proxy proxy_http ssl remoteip
a2enconf bayeos-server.conf
systemctl try-reload-or-restart apache2.service    